# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=gunicorn -b 127.0.0.1:5000 -w 2 --timeout 120 gateway_app:app
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:migrate]
command=bash -c "bash /app/ops/run_marketing_migrations.sh && bash /app/ops/run_producer_migrations.sh && touch /tmp/forge_migrations_done"
directory=/app
environment=DJANGO_SETTINGS_MODULE="logic_service.settings.docker",RUNNING_IN_DOCKER="1"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:marketing]
command=bash -c "i=0; until [ -f /tmp/forge_migrations_done ] || [ $i -ge 120 ]; do sleep 1; i=$((i+1)); done && exec gunicorn -b 127.0.0.1:8002 -w 2 --timeout 120 dashboard.app:app"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:producer]
command=bash -c "if [ \"${ENABLE_PRODUCER:-auto}\" = \"0\" ] || [ ! -d /app/Producer ] || [ ! -f /app/Producer/manage.py ]; then echo 'Producer disabled or missing; skipping Producer process.'; exec tail -f /dev/null; fi; i=0; until [ -f /tmp/forge_migrations_done ] || [ -f /tmp/forge_migrations_failed ] || [ $i -ge 120 ]; do sleep 1; i=$((i+1)); done; if [ -f /tmp/forge_migrations_failed ]; then echo 'ERROR: DB migrations failed — refusing to start Producer to avoid serving a broken schema.'; exit 1; fi; cd /app/Producer; python manage.py loadinitialdata && python manage.py fix_stuck_tasks || true; exec gunicorn -b 127.0.0.1:8001 -w 2 --timeout 1800 logic_service.wsgi"
environment=DJANGO_SETTINGS_MODULE="logic_service.settings.docker",RUNNING_IN_DOCKER="1",FORCE_SCRIPT_NAME="/producer",SCRIPT_NAME="/producer",STATIC_URL="/producer/static/"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
